# services/ppp-connect.timer
[Unit]
Description=ppp connection monitoring

[Timer]
OnBootSec=1m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
